# crontab
*/5 * * * * /usr/bin/python3 /app/run_elt.py >> /var/log/elt.log 2>&1
