# crontab
*/5 * * * * /usr/local/bin/python /app/run_elt.py >> /var/log/elt.log 2>&1
